# systemd/dotpush.service
[Unit]
Description=Backup and push dotfiles using DotPush

[Service]
Type=oneshot
ExecStart=/usr/bin/dotpush backup
ExecStart=/usr/bin/dotpush push
